# etc/cron.d/livefortunes
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/usr/games

9 9,21 * * * samson cd /home/samson/code/livefortunes && ./publish.coffee >> log/publish.log 2>&1
